# runit/log/run
#!/bin/sh
exec chpst -ucoral svlogd -tt ./main
